# pyPit/refs/test.conf

[program:test]
command = sh ./auto_run.sh
directory = /home/ubuntu/test/cm_portal3
user = root
